# example/services.d/service2/finish
#! /bin/sh

echo "service2 is finished"
s6-svscanctl -t /var/run/s6/services
