# example/services.d/service2/run
#!/usr/bin/with-contenv sh

for i in `seq 1 10`
do
	echo "$ECHO_STR$i"
	sleep 5
done
