# example/services.d/service1/finish
#! /bin/sh

echo "service1 is finished"
