# example/services.d/service1/run
#! /bin/sh

exec ping -c 10 google.co.jp
